# base/kernel/files/rc.local
modprobe bridge
modprobe nf_conntrack
modprobe ip_conntrack
modprobe br_netfilter